# deploy/systemd/zephyr-web.service
[Unit]
Description=Zephyr Weather Web UI
Documentation=https://github.com/cngarrison/zephyr

# The web daemon makes server-side HTTP requests to the engine (config.web.engine_url)
# for SSR data fetches, so a working network is required.
After=network-online.target zephyr-engine.service
Wants=network-online.target

# Declare a dependency on the engine.  zephyr-web will still start even if
# zephyr-engine fails (Requires= would stop it); it just won't serve live data
# until the engine comes up.  Use Requires= here if you prefer strict coupling.
BindsTo=zephyr-engine.service

[Service]
Type=exec

# Run as the dedicated zephyr system user.
# Create with: useradd -r -s /sbin/nologin -d /var/lib/zephyr -M zephyr
User=zephyr
Group=zephyr

WorkingDirectory=/var/lib/zephyr

# Tell both daemons where to find the shared TOML config file.
Environment=ZEPHYR_CONFIG=/etc/zephyr/zephyr.toml

# PORT and HOSTNAME are read directly by the Fresh framework before user code
# runs, so they must be set here as process environment variables rather than
# in zephyr.toml.
Environment=PORT=8081
Environment=HOSTNAME=0.0.0.0

ExecStart=/usr/local/bin/zephyr-web

Restart=on-failure
RestartSec=5s

# Shares the /run/zephyr/ runtime directory created by zephyr-engine.service.
RuntimeDirectory=zephyr
RuntimeDirectoryMode=0755

# View logs with:
#   journalctl -u zephyr-web -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zephyr-web

[Install]
WantedBy=zephyr.target
