# deploy/systemd/zephyr-engine.service
[Unit]
Description=Zephyr Weather Engine
Documentation=https://github.com/cngarrison/zephyr

# Wait for a routable network address before starting.
# This is required when storage.provider=mysql and the database is on a remote host.
# If you use SQLite exclusively and want faster boot, you can relax this to
# After=network.target (and remove the Wants= line).
After=network-online.target
Wants=network-online.target

[Service]
Type=exec

# Run as the dedicated zephyr system user (no login shell, no home directory).
# Create with: useradd -r -s /sbin/nologin -d /var/lib/zephyr -M zephyr
User=zephyr
Group=zephyr

# Working directory for the engine process.
WorkingDirectory=/var/lib/zephyr

# Tell both daemons where to find the shared TOML config file.
# The file must exist at this path (or override with ZEPHYR_CONFIG).
Environment=ZEPHYR_CONFIG=/etc/zephyr/zephyr.toml

ExecStart=/usr/local/bin/zephyr-engine

# Restart on abnormal exit (crash, OOM kill, fatal signal).
# Does NOT restart on a clean exit (SIGTERM from `systemctl stop`).
Restart=on-failure
RestartSec=5s

# Create /run/zephyr/ owned by the zephyr user before the service starts.
# Available for Unix sockets, PID files, or other ephemeral runtime artifacts.
# Removed automatically by systemd when the last unit using it stops.
RuntimeDirectory=zephyr
RuntimeDirectoryMode=0755

# Stdout and stderr are captured by the systemd journal automatically.
# View logs with:
#   journalctl -u zephyr-engine -f          # follow live
#   journalctl -u zephyr-engine --since today
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zephyr-engine

[Install]
# zephyr.target pulls in both daemons; WantedBy here means
# `systemctl enable zephyr-engine` registers it under that target.
WantedBy=zephyr.target
